# templates/ultrahook.conf
[program:ultrahook]
command=/usr/local/bin/ultrahook ULTRAHOOK_SUB_DOMAIN 80
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
